# home/retoran/.config/autostart/copyq-server.desktop
[Desktop Entry]
Name=CopyQServer
Exec=copyq --start-server
Type=Application
