# home/retoran/.local/share/systemd/user/thunar-daemon.service
[Unit]
Description="Thunar file manager"

[Service]
ExecStart=/usr/bin/thunar --daemon

[Install]
WantedBy=graphical-session.target
